# webapp.service
[Unit]
Description=systemd
Documentation=sysdtemd
After=network.target

[Service]
WorkingDirectory=/home/csye/webapp/
EnvironmentFile=/home/csye/webapp/.env
Type=simple
User=webappuser
ExecStart=/usr/bin/node /home/admin/webapp/server.js
Restart=on-failure

[Install]
WantedBy=cloud-init.target

